# cadvisor/cadvisor.service
#/usr/lib/systemd/system
[Unit]
Description=cAdvisor

[Service]
ExecStart=/usr/local/bin/cadvisor

[Install]
WantedBy=multi-user.target
